# webby.service
[Unit]
Description=webby
[Service]
ExecStart=/usr/bin/webby
[Install]
WantedBy=multi-user.target